# dot_config/systemd/user/ssh-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ssh-agent -D -a %t/ssh-agent.socket

[Install]
WantedBy=default.target
